# soma.service
[Unit]
Description=Soma — lightweight model-agnostic agent harness
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
# Adapt User/WorkingDirectory + ExecStart Python path to your setup.
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/path/to/soma
ExecStart=/path/to/soma/.venv/bin/python3 /path/to/soma/src/telegram.py
Restart=on-failure
RestartSec=10
# Logs go to journald — view with `journalctl -u soma -f`
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
